# systemd/bybridge.service
[Unit]
Description=ByBridge Daemon — Android Device Bridge
Documentation=https://github.com/yourusername/bybridge-daemon
After=network.target graphical-session.target
Wants=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/bybridge-daemon
Restart=on-failure
RestartSec=3
NotifyAccess=main

# Run as current user (not root by default)
# Override with: systemctl edit bybridge
# User=msfvenom

# Environment for PulseAudio / PipeWire access
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_RUNTIME_PATH=/run/user/1000/pulse
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus

# Needs access to uinput for virtual input
SupplementaryGroups=input

# Data directory
StateDirectory=bybridge
RuntimeDirectory=bybridge

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/bybridge
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=graphical-session.target
